# ROOT/Users/admin/Library/LaunchAgents/com.blackglove.createuser.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.blackglove.createuser</string>
	<key>ProgramArguments</key>
	    <array>
	        <string>/usr/bin/sudo</string>
	        <string>/Library/BlackGlove/fixedUser.sh</string>
	    </array>
	    <key>RunAtLoad</key>
	    <true/>
	    <key>StandardErrorPath</key>
	    <string>/tmp/bg_user.err</string>
	    <key>StandardOutPath</key>
	    <string>/tmp/bs_user.out</string>
	    <key>ThrottleInterval</key>
	    <integer>10</integer>
</dict>
</plist>
